# supervisor/require_itself.conf
[program:require_itself]
directory=/opt/microservice/src/local_magellan
command=python require_service.py 2000 example-udp
