# supervisor/show_env.conf
[program:show_env]
directory=/opt/example-multi/src
command=python show_env.py 80

[program:register_ssh]
directory=/opt/microservice/src
command=python register_in_service_discovery.py 22 -s ssh
